# site/minecraft/files/minecraft.service.epp
[Unit]
Description=Minecraft Server
Wants=network.target
After=network-online.target
Requires=network-online.target

[Service]
WorkingDirectory=<% $install_dir %>
ExecStart=/usr/bin/java -Xmx512M -Xms32M -jar server.jar nogui
Restart=always

[Install]
WantedBy=multi-user.target
